# docker/app/rootfs/etc/service/nginx/run
#!/bin/bash
set -e

NGINX_CONF="/etc/nginx/nginx.conf"

/usr/sbin/nginx -c $NGINX_CONF -t && exec nginx -c $NGINX_CONF
