# docker/app/rootfs/etc/service/php-fpm/run
#!/bin/sh

PHP_INI="/etc/php/7.0/fpm/php.ini"
PHPFPM_CONF="/etc/php/7.0/fpm/php-fpm.conf"

/usr/sbin/php-fpm7.0 -c $PHP_INI -y $PHPFPM_CONF
